# extra/supervisord/zpark_celery.ini
[program:zpark_celery]
command=/home/zpark/venvs/zpark/bin/celery worker -A zpark:celery
directory=/home/zpark/zpark
user=_zpark
group=_zpark
numprocs=1
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
